# serviceFiles/spothintamqtt.target
[Unit]
Description=spothintamqtt service
StopWhenUnneeded=yes